# packer/_shared/ansible/roles/vitess_firstboot/templates/vitess-node-firstboot.service.j2
[Unit]
Description=vitess-node firstboot (NIC discovery, hostname mapping, VMnet10 backplane, node-identity env)
After=systemd-networkd.service network-online.target
Wants=network-online.target
# No Before= for the role services -- firstboot does NOT enable any
# etcd/vttablet/vtgate/mysqld service. The Terraform role-overlays render
# per-host config (topo addresses, tablet alias, cert SANs) then bring up
# the right component per node post-apply.
ConditionPathExists=!/var/lib/vitess-node-firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vitess-node-firstboot.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
